# pi-s3-sync/systemd/sd-s3-sync.service
[Unit]
Description=Sync S3 bucket to SD card or HOME when wifi.json is present
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/epaper-settings/s3-sync.env
Type=oneshot
User=%EPAPER_USER%
Group=%EPAPER_USER%
ExecStart=/usr/local/bin/sync_s3_from_sd.py

# --- Sandboxing ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
PrivateDevices=true

# SD card read/write
ReadWritePaths=/mnt/epaper_sd

[Install]
WantedBy=multi-user.target
